# essential-files/etc/systemd/system/user@.service.d/dbus.conf
[Service]
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%I/dbus/user_bus_socket
Environment=QT_QPA_PLATFORM=wayland-egl
Environment=WAYLAND_DISPLAY=wayland-%I
